# packaging/phorum.socket
[Unit]
Description=Socket for the Phorum Daemon

[Socket]
BindIPv6Only=both
ListenStream=[::]:70